# crontab
@daily rm /home/arunscape/Downloads/*.{mkv,mp4,webm}
@daily rsync -r -a --delete --exclude=node_modules --exclude=target --exclude=build --exclude=.venv  ~/code/. /run/media/arunscape/znas/dev -P
