# .config/rc/init.d/deemix
#!/sbin/openrc-run

description="Arun's $RC_SVCNAME instance with podman"
run_path="${XDG_RUNTIME_DIR}/podman"
log_path="${run_path}/log"
pidfile="${run_path}/${RC_SVCNAME}.pid"
output_log="${log_path}/${RC_SVCNAME}.log"
error_log=$output_log
command=podman
command_args="run
            --pidfile $pidfile
            --name $RC_SVCNAME
            -v $config_path:/config
            -v $data_path:/data
            -v $downloads_path:/downloads 
            --network $network
            --replace
            -p $web_port:6595
            -e DEEMIX_MUSIC_DIR=/data/music
            -e PUID=0
            -e PGID=0
            -e UMASK_SET=022
            $container_url
            "
command_background=true
supervisor=supervise-daemon

depend() {
    use podman
    need gluetun
}

start_pre(){
  checkpath -o $USER -d "${run_path}" "${log_path}"

  ebegin "Pulling $container_url"
  podman pull $container_url
  eend $? "Failed to pull $container_url"
}
