# .config/rc/init.d/forgejo
#!/sbin/openrc-run

description="Arun's $RC_SVCNAME instance with podman"
run_path="${XDG_RUNTIME_DIR}/podman"
log_path="${run_path}/log"
pidfile="${run_path}/${RC_SVCNAME}.pid"
output_log="${log_path}/${RC_SVCNAME}.log"
error_log=$output_log
command=podman
command_args="run
              --pidfile=$pidfile
              --name $RC_SVCNAME
              -p $web_port:3000
              -p $ssh_port:2222
              -v $data_path:/var/lib/gitea
              -v $conf_path:/etc/gitea
              --userns=keep-id
              --cgroup-manager=cgroupfs
              --replace
              $container_url
             "
command_background=true


depend() {
    use podman
}

start_pre(){
  checkpath -o $USER -d "${run_path}" "${log_path}"

  ebegin "Checking $data_path and $conf_path"
  checkpath -o $USER -d "${data_path}" "${conf_path}"
  eend $? "Double check $data_path and $conf_path"

  ebegin "Pulling $container_url"
  podman pull $container_url
  eend $? "Failed to pull $container_url"
}
