# .config/rc/init.d/ollama
#!/sbin/openrc-run

description="Arun's $RC_SVCNAME instance with podman"
run_path="${XDG_RUNTIME_DIR}/podman"
log_path="${run_path}/log"
pidfile="${run_path}/${RC_SVCNAME}.pid"
output_log="${log_path}/${RC_SVCNAME}.log"
error_log=$output_log
command=podman
command_args="run
              --pidfile $pidfile
              --name $RC_SVCNAME
              --device /dev/kfd
              --device /dev/dri
              -v $data_path:/root/.ollama
              -p $web_port:11434
              --replace
              --network $network
              $container_url
             "
command_background=true
supervisor=supervise-daemon

depend() {
    use podman
}

start_pre(){
  checkpath -o $USER -d "${run_path}" "${log_path}"

  ebegin "Checking $data_path"
  checkpath -o $USER -d "${data_path}"
  eend $? "Double check $data_path"

  ebegin "Pulling $container_url"
  podman pull $container_url
  eend $? "Failed to pull $container_url"
}

stop_post(){
  ebegin "Stopping $RC_SVCNAME"
  podman stop $RC_SVCNAME
  eend $? "Failed to stop $RC_SVCNAME"
}
