# .config/rc/init.d/sonarr
#!/sbin/openrc-run

description="Arun's $RC_SVCNAME instance with podman"
run_path="${XDG_RUNTIME_DIR}/podman"
log_path="${run_path}/log"
pidfile="${run_path}/${RC_SVCNAME}.pid"
output_log="${log_path}/${RC_SVCNAME}.log"
error_log=$output_log
command=podman
command_args="run
              --pidfile $pidfile
              --name $RC_SVCNAME
              -e PUID=0
              -e PGID=0
              -e UMASK=002
              -e TZ=$tz
              -p $web_port:8989
              -v $config_path:/config
              -v $data_path:/data
              --replace
              --network $network
              $container_url
             "
command_background=true
supervisor=supervise-daemon

depend() {
    use podman
    need prowlarr
}

start_pre(){
  checkpath -o $USER -d "${run_path}" "${log_path}"

  ebegin "Pulling $container_url"
  podman pull $container_url
  eend $? "Failed to pull $container_url"
}

stop_post(){
  ebegin "Stopping $RC_SVCNAME"
  podman stop $RC_SVCNAME
  eend $? "Failed to stop $RC_SVCNAME"
}
