# .config/systemd/whoami.socket
[Socket]
ListenStream = 1234
Service = whoami.service

[Install]
WantedBy = sockets.target
